# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/var/run/supervisord.pid

# ── PostgreSQL ────────────────────────────────────────────────────────────────
[program:postgres]
# Run postgres in foreground with Debian's config file location
command=/usr/lib/postgresql/15/bin/postgres
    -D /var/lib/postgresql/15/main
    --config-file=/etc/postgresql/15/main/postgresql.conf
    -c shared_buffers=16MB
    -c max_connections=25
    -c work_mem=1MB
    -c maintenance_work_mem=4MB
    -c effective_cache_size=32MB
    -c wal_buffers=1MB
    -c checkpoint_completion_target=0.9
    -c autovacuum=off
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/trackshop/postgres.log
stderr_logfile=/var/log/trackshop/postgres.log
startsecs=3
# ── dnsmasq DNS service ─────────────────────────────────────────────────────────
# Intentionally vulnerable: cache-size=0, no rate limit, oversized TXT records
[program:dnsmasq]
command=dnsmasq --no-daemon --conf-file=/etc/dnsmasq.d/lab.conf --log-queries --log-facility=-
autostart=true
autorestart=true
priority=15
stdout_logfile=/var/log/trackshop/dnsmasq.log
stderr_logfile=/var/log/trackshop/dnsmasq.log
startsecs=2
# ── Node.js backend API ───────────────────────────────────────────────────────
[program:backend]
# Cap the V8 heap to keep Node.js memory footprint small
command=/usr/bin/node --max-old-space-size=96 /app/server.js
directory=/app
environment=
    DATABASE_URL="postgres://appuser:apppassword@127.0.0.1:5432/trackshop",
    NODE_ENV="production",
    PORT="3000"
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/trackshop/backend.log
stderr_logfile=/var/log/trackshop/backend.log
; Give postgres time to be ready before counting backend as stable
startsecs=5

# ── Nginx ─────────────────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/trackshop/nginx.log
stderr_logfile=/var/log/trackshop/nginx.log
startsecs=2
